# nhk-news.service
[Unit]
Description=NHK News API Service
After=network.target

[Service]
Type=simple
User=tetsuya
WorkingDirectory=/Users/tetsuya/Dev/Fudoki4Web
Environment=PATH=/Users/tetsuya/Dev/Fudoki4Web/python_scripts/venv/bin
ExecStart=/Users/tetsuya/Dev/Fudoki4Web/python_scripts/venv/bin/python /Users/tetsuya/Dev/Fudoki4Web/python_scripts/nhk_rss_scraper.py api
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
